# debian/automounter.service
[Unit]
Description=A daemon for auto-mounting USB devices.
Wants=usbmuxd.service

[Service]
ExecStart=/usr/bin/automounter